# system_files/usr/lib/systemd/user/brave-unlock.service
[Unit]
Description=Remove stale Brave profile singleton lock at login
After=default.target

[Service]
Type=oneshot
ExecStart=/usr/bin/rm -f %h/.config/BraveSoftware/Brave-Browser/SingletonLock %h/.config/BraveSoftware/Brave-Browser/SingletonCookie %h/.config/BraveSoftware/Brave-Browser/SingletonSocket

[Install]
WantedBy=default.target
