# system_files/usr/lib/systemd/user/nextcloud-unlock.service
[Unit]
Description=Remove stale Nextcloud kdsingleapplication lock at login
After=default.target

[Service]
Type=oneshot
ExecStart=/usr/bin/rm -f %h/.var/app/com.nextcloud.desktopclient.nextcloud/cache/tmp/kdsingleapp-%u-nextcloud %h/.var/app/com.nextcloud.desktopclient.nextcloud/cache/tmp/kdsingleapp-%u-nextcloud.lock

[Install]
WantedBy=default.target
