# system_files/usr/lib/systemd/system/bazzite-custom-flatpaks.service
[Unit]
Description=Apply bazzite-custom flatpak preinstall list (Fedora-native mechanism)
# Bazzite's flatpak-manager configures Flathub as the system remote first;
# we wait for that and for network before pulling app metadata.
After=bazzite-flatpak-manager.service network-online.target
Wants=network-online.target
Requires=bazzite-flatpak-manager.service
# Don't pin to a specific image — we just need to be on a deploy that ships
# the .preinstall file at /usr/share/flatpak/preinstall.d/.
ConditionPathExists=/usr/share/flatpak/preinstall.d/bazzite-custom.preinstall

[Service]
Type=oneshot
RemainAfterExit=yes
# `flatpak preinstall` is idempotent and respects user uninstalls (per
# `man flatpak-preinstall`). Re-runs every boot are cheap if nothing
# changed. -y autoselects defaults; --noninteractive skips prompts.
ExecStart=/usr/bin/flatpak preinstall -y --noninteractive
# Network/Flathub hiccups on boot — let systemd retry the unit.
Restart=on-failure
RestartSec=60
StartLimitInterval=0

[Install]
WantedBy=multi-user.target
